# octoprint_mrbeam/files/migrate/mount_manager_add.service
[Service]
Type=oneshot
ExecStart=/usr/bin/mount_manager add

[Unit]
After=mount_manager_remove.service
