# octoprint_mrbeam/files/migrate/mount_manager_remove.service
[Service]
Type=oneshot
ExecStart=/usr/bin/mount_manager remove
